# scripts/led-matrix/service_files/led-matrix-web.service
[Unit]
Description=LED Matrix Web Controller
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/led_matrix
ExecStart=/usr/bin/python3 /opt/led_matrix/led_matrix_web.py --config /etc/led-matrix/led_matrix_display_config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
